# themes/cambium/root/etc/uci-defaults/luci-theme-cambium
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Cambium=/luci-static/cambium
        commit luci
EOF
	
